# sensdle.service
[Unit]
Description=Sensdle Port Forward Service
Documentation=https://github.com/your-repo/sensdle
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=matty
Group=matty
WorkingDirectory=/home/matty/terraform/monitoring2/nodestuff/ottawa-senators-daily
ExecStart=/home/matty/terraform/monitoring2/nodestuff/ottawa-senators-daily/scripts/port-forward.sh
Restart=always
RestartSec=5
StartLimitInterval=0

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/matty/terraform/monitoring2/nodestuff/ottawa-senators-daily

# Environment
Environment=HOME=/home/matty
Environment=KUBECONFIG=/home/matty/.kube/config

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sensdle

[Install]
WantedBy=multi-user.target